# config/web.supervisord.conf.erb
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)

[program:nginx]
command=bin/nginx -p . -c config/nginx.conf
redirect_stderr=true
autorestart=true
stdout_logfile=logs/access.log
stderr_logfile=logs/error.log

[program:web]
command=python -m app.web 888%(process_num)d
process_name=%(program_name)s_%(process_num)d
redirect_stderr=true 
stdout_logfile=logs/access.log
stderr_logfile=logs/error.log
autorestart=true
numprocs=<%= ENV['WEB_PROCESSES'] || 2 %>
;numprocs_start=100
priority=1
