# config/worker.supervisord.conf.erb
[supervisord]
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)

[program:worker]
command=python -m app.worker
process_name=%(program_name)s_%(process_num)d
redirect_stderr=true 
stdout_logfile=logs/access.log
stderr_logfile=logs/error.log
autorestart=true
numprocs=<%= ENV['WORKER_PROCESSES'] || 4 %>
